# utils/monglue.crontab
#!/bin/sh

# @hourly bash /home/ocw/ebill2/utils/monglue.crontab

unset http_proxy
unset https_proxy

ruby /var/lib/ebill-monglue/monglue-clean.rb > /tmp/monglue-clean.log 2>&1
